# launchd/com.marcilio.daily-agentic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marcilio.daily-agentic</string>

    <!-- Invoke the wrapper: agent run → build dist → git push (Vercel deploys) -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/marcilio/projects/researcher-agent-dashboard/scripts/daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/marcilio/projects/researcher-agent-dashboard</string>

    <!-- Run twice weekly: Monday and Thursday at 02:00 local time
         (Weekday: 0=Sun, 1=Mon, 2=Tue, 3=Wed, 4=Thu, ..., 6=Sat) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Logs: tail -f either of these to watch the daily cycle -->
    <key>StandardOutPath</key>
    <string>/Users/marcilio/Library/Logs/daily-agentic.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/marcilio/Library/Logs/daily-agentic.err.log</string>

    <!-- Don't fire immediately on load; wait for the next 05:00 -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
